# infra/gateway.service
# RPM Gateway Service — systemd unit file
# Install on the Raspberry Pi:
#
#   sudo cp gateway.service /etc/systemd/system/rpm-gateway.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable rpm-gateway
#   sudo systemctl start rpm-gateway
#   sudo journalctl -u rpm-gateway -f
#
# Edit Environment lines or use an EnvironmentFile for secrets.

[Unit]
Description=RPM Gateway — BLE sensor ingest and cloud sync
After=network-online.target bluetooth.target
Wants=network-online.target bluetooth.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/mdmcg

# Python virtual environment
ExecStart=/home/pi/mdmcg/.venv/bin/python -m gateway

# Secrets and config — stored in env file, never committed to git
# Copy infra/.rpm_gateway.env.example -> /home/pi/.rpm_gateway.env and fill in values
EnvironmentFile=/home/pi/.rpm_gateway.env

# Restart policy: always restart on failure with a delay
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rpm-gateway

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pi/
PrivateTmp=true

[Install]
WantedBy=multi-user.target
